# services.d/gor/run
#!/usr/bin/with-contenv sh

/gopath/bin/gor -input-raw :80 -output-http=$GOR_OUTPUT_HTTP